# cloud-init.init
#!/sbin/openrc-run

name="cloud-init"
command="/bin/cloud-init"

depend() {
	before sshd
	need localmount
	need logger
	need networking
}

